# src/etc/xdg/autostart/qubes-vpn-notifier.desktop.in
[Desktop Entry]
Version=1.0
Encoding=UTF-8
Name=Qubes VPN notifications
Exec=@LIBEXECDIR@/qubes-vpn-notifier
Terminal=False
Type=Application
Categories=
GenericName=
X-GNOME-Autostart-Phase=Initialization
OnlyShowIn=X-QUBES;
