# src/usr/lib/systemd/system/qubes-vpn-configuration.path.in
[Unit]
Description=Qubes OS leakproof VPN configuration monitor
ConditionPathExists=@QUBESSERVICEDIR@/qubes-vpn

[Path]
PathChanged=@VPNCONFDIR@/@VPNCONFFILE@
PathChanged=@VPNCONFDIR@

[Install]
WantedBy=network.target
