# src/usr/lib/systemd/system/qubes-vpn.service.in
[Unit]
Description=Qubes OS leakproof VPN
ConditionPathExists=@VPNCONFDIR@/@VPNCONFFILE@
ConditionPathExists=@QUBESSERVICEDIR@/qubes-vpn
After=network.target network.target qubes-network.service qubes-vpn-forwarding.service

[Service]
PrivateTmp=true
Type=forking
PIDFile=@VPNRUNDIR@/pid
ExecStartPre=@BINDIR@/mkdir -p @VPNRUNDIR@
ExecStartPre=@BINDIR@/chown -R qubes-vpn @VPNRUNDIR@
ExecStartPre=@BINDIR@/chgrp -R qubes-vpn @VPNCONFDIR@
ExecStartPre=@BINDIR@/chmod -R g+rX @VPNCONFDIR@
ExecStart=@SBINDIR@/openvpn --daemon --writepid @VPNRUNDIR@/pid --cd @VPNCONFDIR@ --config @VPNCONFFILE@ --route-noexec --script-security 2 --up @SBINDIR@/qubes-vpn-interface-control --down @SBINDIR@/qubes-vpn-interface-control --up-restart --user qubes-vpn --group qubes-vpn --ifconfig-noexec --persist-tun --mlock
ExecStartPost=@BINDIR@/chown qubes-vpn @VPNRUNDIR@/pid
Restart=on-failure

[Install]
WantedBy=multi-user.target
