# src/usr/lib/systemd/system/qubes-vpn-configuration.service.in
[Unit]
Description=Qubes OS leakproof VPN configuration restarter
ConditionPathExists=@QUBESSERVICEDIR@/qubes-vpn

[Service]
Type=oneshot
ExecStart=@BINDIR@/qubes-vpn-config-change-detector
