# src/usr/lib/systemd/system/qubes-vpn-configuration.service.in
[Unit]
Description=Qubes OS leakproof VPN configuration restarter
ConditionPathExists=@QUBESSERVICEDIR@/qubes-vpn

[Service]
Type=oneshot
ExecStartPre=@BINDIR@/sleep 0.5
ExecStart=@BINDIR@/systemctl --no-block restart qubes-vpn.service
ExecStartPost=@BINDIR@/sleep 0.5
